# .config/systemd/user/restic-backup-home.timer
[Unit]
Description=Daily backup home
[Timer]
OnCalendar=daily
Persistent=true
[Install]
WantedBy=timers.target
