# .config/systemd/user/restic-prune.service
[Unit]
Description=Prune restic backups
[Service]
Type=oneshot
ExecStart=/snap/bin/restic prune
EnvironmentFile=%h/.config/restic/secrets
