# .config/systemd/user/restic-prune.service
[Unit]
Description=Prune restic backups
[Service]
Type=oneshot
ExecStart=restic prune
EnvironmentFile=%h/.config/restic/secrets
