# .config/autostart/gpaste.desktop
[Desktop Entry]
Type=Application
Name=GPaste
Icon=edit-paste
Description=
Comment=Clipboard manager
Exec=gpaste-client start
